# ansible/duplicacy/duplicacy.timer
[Unit]
Name=Duplicacy Backup

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target