# ansible/roles/snapraid/files/snapraid.service
[Unit]
Description=SnapRAID sync and scrub
OnFailure=status-email-user@%n.service

[Service]
ExecStart=/bin/bash -c "snapraid sync && snapraid scrub"