# ansible/roles/duplicacy/files/duplicacy.service
[Unit]
Name=Duplicacy Backup

[Service]
ExecStart=/usr/local/bin/duplicacy backup
WorkingDirectory=/mnt/storage/Kevin